# music-bot.service
[Unit]
Description=Music Discord Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=${BOT_USER}
WorkingDirectory=${BOT_DIR}
ExecStart=${BOT_PYTHON} bot.py
EnvironmentFile=${BOT_DIR}/.env

Restart=on-failure
RestartSec=10
StartLimitBurst=5
StartLimitInterval=300

TimeoutStartSec=60
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

StandardOutput=journal
StandardError=journal
SyslogIdentifier=music-bot

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=${BOT_DIR}

[Install]
WantedBy=multi-user.target
